# httpd/container/usr/lib/systemd/system/httpd.service.d/environment.conf
[Service]
ExecStartPre=/usr/bin/bash -c "/usr/local/bin/generate-plain-environment-file"
EnvironmentFile=/etc/profile.d/environment.env
Restart=always
RestartSec=3
